# automation/com.codex.nmpa-toothpaste-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.nmpa-toothpaste-tracker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/liuaiwei/Downloads/AI/Codex/国产牙膏备案信息爬取/automation/run_monthly_tracker.sh</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/liuaiwei/Downloads/AI/Codex/国产牙膏备案信息爬取</string>

  <key>StandardOutPath</key>
  <string>/Users/liuaiwei/Downloads/AI/Codex/国产牙膏备案信息爬取/output/logs/monthly-task.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/liuaiwei/Downloads/AI/Codex/国产牙膏备案信息爬取/output/logs/monthly-task.stderr.log</string>
</dict>
</plist>
